# files/profile.d/custom_aliases.sh
alias tmux="/usr/local/bin/tmux.sh"
alias glances="glances --disable-log --enable-process-extended"
